# Task 1/crontab
# Run ingestion script daily at 1:00 AM
0 1 * * * python3 /app/ingestions.py --single-day >> /var/log/ingestion/ingestion.log 2>&1